# ng-backend/ai-monitor.service
[Unit]
Description=Moomoo Bot AI Monitor Analysis Service
After=network.target

[Service]
Type=oneshot
User=clck
WorkingDirectory=/home/clck/moomoo-bot/ng-backend
ExecStart=/home/clck/moomoo-bot/ng-backend/venv/bin/python ai-monitor/run_analysis.py
StandardOutput=append:/home/clck/moomoo-bot/ng-backend/ai-monitor.log
StandardError=append:/home/clck/moomoo-bot/ng-backend/ai-monitor.err
EnvironmentFile=/home/clck/moomoo-bot/ng-backend/.env

[Install]
WantedBy=multi-user.target
