# ng-backend/ai-monitor.timer
[Unit]
Description=Run AI Monitor Analysis Daily at 07:00 and 18:00 US/Eastern

[Timer]
# Schedule: 07:00 and 18:00 (Local Time of VM)
OnCalendar=*-*-* 07:00:00
OnCalendar=*-*-* 18:00:00
Persistent=true
Unit=ai-monitor.service

[Install]
WantedBy=timers.target
